# packages/operbot/operbot-200-py3-none-any.whl/operbot-200.data/data/operbot/operbot.service
[Unit]
Description=operator bot
After=multi-user.target

[Service]
Type=forking
KillMode=mixed
User=operbot
Group=operbot
DynamicUser=yes
StateDirectory=operbot
LogsDirectory=operbot
CacheDirectory=operbot
ExecStart=/usr/local/bin/operbotd
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
