# packages/enough/enough-2.2.35.tar.gz/enough-2.2.35/playbooks/forgejo/roles/forgejo-runner/templates/forgejorunner.service.j2
[Unit]
Description=Forgejo runner
After=syslog.target network.target

[Service]
Type=simple
User=debian
Group=debian
WorkingDirectory=/srv/forgejorunner
ExecStart=/bin/forgejo-runner daemon

[Install]
WantedBy=multi-user.target
