# packages/rssbot/rssbot-330-py3-none-any.whl/rssbot-330.data/data/rssbot/rssbot.service
[Unit]
Description=RSSBOT - feeding rss into your channel
After=multi-user.target

[Service]
#Type=simple
User=rssbot
Group=rssbot
DynamicUser=yes
StateDirectory=rssbot
LogsDirectory=rssbot
CacheDirectory=rssbot
ExecStart=/usr/local/bin/rssbotd
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
