# packages/genocide/genocide-111.tar.gz/genocide-111/files/genocide.service
[Unit]
Description=OTP-CR-117/19
After=multi-user.target

[Service]
Type=forking
KillMode=mixed
User=genocide
Group=genocide
DynamicUser=yes
StateDirectory=genocide
LogsDirectory=genocide
CacheDirectory=genocide
ExecStart=/usr/local/bin/genocided
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
